# .config/systemd/user/battery-alert.service
[Unit]
Description=Desktop alert warning of low remaining battery

[Service]
Type=oneshot
ExecStart=%h/.local/bin/battery-alert

[Install]
WantedBy=graphical.target
